# packages/devbox/etc/cron.daily/manual_packages_backup
#! /bin/sh
#  by: Andrew Velez

_logged_err_message="Error running xpkg -m > manual_packages"
_pkgs_backup_file="/home/andrew/Code/.devbox/manual_packages"

[ -e "${_pkgs_backup_file}" ] || touch "${_pkgs_backup_file}"

if xpkg -m > "${_pkgs_backup_file}"; then
  vlogger -p cron.err -t manual_packages_backup "${_logged_err_message}"
  exit 1
fi
